# configs/99-usb-auto.rules
ACTION=="add", SUBSYSTEM=="block", KERNEL=="sd[a-z][0-9]", ENV{ID_DRIVE_THUMB}=="1", RUN+="/usr/local/bin/usb-auto.sh %N"
