# configs/usb-auto@.service
[Unit]
Description=USB auto handler for %I

[Service]
Type=oneshot
ExecStart=/usr/local/bin/usb-auto.sh %I
